# .config/systemd/user/hyde-config.service
[Unit]
Description=HyDE Configuration Parser Service
Documentation=https://github.com/HyDE-Project/hyde-config
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/lib/hyde/hyde-config
Restart=on-failure
RestartSec=5s
Environment="DISPLAY=:0"

# Make sure the required directories exist
ExecStartPre=/bin/mkdir -p %h/.config/hyde
ExecStartPre=/bin/mkdir -p %h/.local/state/hyde

[Install]
WantedBy=graphical-session.target